# supervisor-supercronic.conf
[program:supercronic]
command=/usr/local/bin/supercronic -split-logs -no-reap /etc/crontab.supercronic
redirect_stderr=true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
autostart = true
autorestart = true
priority = 20
